# sv/bhyve-openbsd/finish
#!/bin/sh
set -eu
[ -r conf ] && . ./conf
: ${BHYVECTL:=/usr/sbin/bhyvectl}
name=${PWD##*-}
[ -r /dev/vmm/${name} ] && ${BHYVECTL} --vm=${name} --destroy
for intf in $(ifconfig -g bhyve-${name}); do
	ifconfig ${intf} destroy
done
