# sv/jackd/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${USER:=root}"
: "${OPTS:=-r -doss -r44100 -p1024 -n3 -w16}"
# shellcheck disable=SC2086
exec chpst -u ${USER} /usr/local/bin/jackd ${OPTS}
