# etc/sv/ifstated-template/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${OPTS:=-f ${PWD}/ifstated.conf}"
# shellcheck disable=SC2086
exec /usr/local/sbin/ifstated -d ${OPTS}
