# etc/sv/zfsd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/sbin/zfsd -d
