# etc/sv/newsyslog/run
#!/bin/sh
[ -r conf ] && . ./conf
# Rotate log files every hour, if necessary.
exec snooze -H '*' \
	env PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/sbin \
	/usr/sbin/newsyslog
