# etc/sv/php-fpm/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/local/sbin/php-fpm --nodaemonize
