# sv/nmbd/run
#!/bin/sh
[ -r conf ] && . ./conf
mkdir -p /var/run/samba4
exec nmbd -F -S ${OPTS}
