# sv/opensmtpd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/local/sbin/smtpd -d ${OPTS} 2>&1 > /dev/null
