# sv/opensmtpd/run
#!/bin/sh
[ -r conf ] && . ./conf
# shellcheck disable=SC2086
exec /usr/local/sbin/smtpd -d ${OPTS} > /dev/null 2>&1
