# sv/jail-template/finish
#!/bin/sh
[ -r conf ] && . ./conf
name=${PWD##*-}
: ${ROOT:=/usr/jails/${name}}
: ${CONFIG:=${PWD}/jail.conf}
jail -f ${CONFIG} ${OPTS} -r ${name} 2>&1
umount -f ${ROOT}/var/run
