# sv/webcamd-template/run
#!/bin/sh
[ -r conf ] && . ./conf
kldload -qn cuse
# shellcheck disable=SC2086
exec /usr/local/sbin/webcamd -N "${PWD#*@}" ${OPTS} 2>&1
