# etc/sv/lpd/run
#!/bin/sh
set -e
[ -r conf ] && . ./conf
: "${OPTS:=-s}"
chkprintcap -d
# shellcheck disable=SC2086
exec lpd -f ${OPTS}
