# sv/webcamd-sample/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/local/sbin/webcamd -N ${PWD#*-} ${OPTS} 2>&1
