# sv/runsvdir-user-template/finish
#!/bin/sh
# shellcheck disable=SC2086
USER=${PWD##*@}
svdir=$(eval echo "~${USER%%:*}/service")
sv force-stop ${svdir}/*
sv exit ${svdir}/*
