# sv/runsvdir-user-template/finish
#!/bin/sh
USER=${PWD##*@}
svdir=$(eval echo ~${USER%%:*}/service)
sv force-stop ${svdir}/*
sv exit ${svdir}/*
