# sv/cupsd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/local/sbin/cupsd -f ${OPTS}
