# sv/lpd/run
#!/bin/sh
set -e
[ -r conf ] && . ./conf
: ${OPTS:=-s}
chkprintcap -d
exec lpd -f ${OPTS}
