# sv/rtsold/run
#!/bin/sh
[ -r conf ] && . ./conf
: ${OPTS:=-a}
exec /usr/sbin/rtsold -f ${OPTS}
