# sv/autounmountd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/sbin/autounmountd -d ${OPTS} 2>&1
