# etc/sv/kore-template/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${CONFIG:=/dev/null}"
: "${KORE:=/usr/local/bin/kore}"
# shellcheck disable=SC2086
exec ${KORE} -f -c ${CONFIG} ${OPTS} 2>&1
