# sv/svnserve/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${USER:=svn:svn}"
# shellcheck disable=SC2086
exec chpst -u ${USER} /usr/local/bin/svnserve --foreground -d ${OPTS}
