# sv/devd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /sbin/devd -d ${OPTS} > /dev/null 2>&1
