# etc/sv/openntpd/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${OPTS:=-s}"
# shellcheck disable=SC2086
exec /usr/local/sbin/ntpd -d ${OPTS} 2>&1
