# sv/ifstated-template/run
#!/bin/sh
[ -r conf ] && . ./conf
: ${OPTS:=-f ${PWD}/ifstated.conf}
exec /usr/local/sbin/ifstated -d ${OPTS}
