# sv/namecheap-ddns/run
#!/bin/sh
exec snooze -M '/5' -H '*' ./namecheap-ddns
