# sv/dhclient-sample/run
#!/bin/sh
[ -r conf ] && . ./conf
intf=${PWD##*-}
# Wait for interface
while ! ifconfig ${intf} > /dev/null 2>&1; do
        sleep 1
done
exec dhclient -d ${intf} 2>&1
