# sv/jail-sample/finish
#!/bin/sh
[ -r conf ] && . ./conf
name=${PWD##*-}
exec jail ${OPTS} -r ${name} 2>&1
