# sv/musicpd/run
#!/bin/sh
[ -r conf ] && . ./conf
: ${OPTS:=/usr/local/etc/musicpd.conf}
exec /usr/local/bin/musicpd --no-daemon ${OPTS}
