# sv/musicpd/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${OPTS:=/usr/local/etc/musicpd.conf}"
# shellcheck disable=SC2086
exec /usr/local/bin/musicpd --no-daemon ${OPTS}
