# sv/acme-client@example.com@www.example.com/run
#!/bin/sh
exec snooze -R 600 -H1 ./acme-client.sh
