# sv/runsvdir-user/run
#!/bin/sh
USER=${PWD##*-}
exec 2>&1
install -d -m 700 -o ${USER%%:*} "$(eval echo ~${USER%%:*})/.runit-supervise"
exec chpst -u ${USER} -P runsvdir -P "$(eval echo ~${USER%%:*}/service)" 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'
