# sv/runsvdir-user/run
#!/bin/sh
USER=${PWD##*-}
exec 2>&1
exec chpst -u ${USER} -P runsvdir -P $(eval echo ~${USER%%:*}/service) 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'
