# sv/netdata/run
#!/bin/sh
[ -r conf ] && . ./conf
: ${USER:=netdata:netdata}
exec chpst -u ${USER} /usr/local/sbin/netdata -D ${OPTS}
