# etc/sv/getty-ttyv1/run
#!/bin/sh
tty=${PWD##*-}
[ -r conf ] && . ./conf
: "${TYPE:=Pc}"
: "${TERM:=xterm}"
[ ! -c /dev/${tty} ] && sv d getty-${tty}
export TERM
exec /usr/libexec/getty ${TYPE} "${tty}"
