# sv/obhttpd/run
#!/bin/sh
[ -r conf ] && . ./conf
# shellcheck disable=SC2086
exec /usr/local/sbin/obhttpd -d ${OPTS}
