# etc/sv/nmbd/run
#!/bin/sh
[ -r conf ] && . ./conf
mkdir -p /var/run/samba4
# shellcheck disable=SC2086
exec nmbd -F -S ${OPTS}
