# sv/hiawatha/run
#!/bin/sh
[ -r conf ] && . ./conf
exec 2>&1
# shellcheck disable=SC2086
exec hiawatha -d ${OPTS} >/dev/null
