# sv/kodi-gbm/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${USER:=kodi:kodi}"
exec chpst -u ${USER}:video env HOME="$(eval echo "~${USER%%:*}")" /usr/local/bin/kodi 2>&1
