# etc/sv/rtsold/run
#!/bin/sh
[ -r conf ] && . ./conf
: "${OPTS:=-a}"
# shellcheck disable=SC2086
exec /usr/sbin/rtsold -f ${OPTS}
