# sv/dbus/run
#!/bin/sh
[ -r conf ] && . ./conf
/usr/local/bin/dbus-uuidgen --ensure
mkdir -p /var/run/dbus
# shellcheck disable=SC2086
exec /usr/local/bin/dbus-daemon --system --nofork --nopidfile ${OPTS}
