# sv/zerotier/run
#!/bin/sh
kldload -qn if_tuntap
kldload -qn if_tap
[ -r conf ] && . ./conf
# shellcheck disable=SC2086
exec /usr/local/sbin/zerotier-one ${OPTS}
