# etc/sv/nginx/run
#!/bin/sh
[ -r conf ] && . ./conf
# shellcheck disable=SC2086
exec nginx -g 'daemon off;' ${OPTS}
