# sv/cron/run
#!/bin/sh
[ -r conf ] && . ./conf
exec /usr/sbin/cron -n $OPTS 2>&1
