# etc/sv/automountd/run
#!/bin/sh
[ -r conf ] && . ./conf
/usr/sbin/automount
# shellcheck disable=SC2086
exec /usr/sbin/automountd -d ${OPTS} >/dev/null 2>&1
