# sv/periodic-monthly/run
#!/bin/sh
# Perform monthly maintenance.
exec snooze -M 30 -H 5 -d 1 \
	env \
	HOME=/root \
	PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/sbin \
	/usr/sbin/periodic monthly 2>&1
