# sv/powerdxx/run
#!/bin/sh
[ -r conf ] && . ./conf
# shellcheck disable=SC2086
exec /usr/local/sbin/powerd++ -f ${OPTS} > /dev/null 2>&1
